# assignement_4_MP3_player/myoverlay/etc/profile.d/playmusic.sh
#!/bin/sh
export PS1='<rapi_shell> : \w\$'
modprobe snd-bcm2835
. /myApplications/buttonmanager.sh 4 17 27 22 &

